# tests/cases/noimage.container
## assert-failed
## assert-stderr-contains "Converting \"noimage.container\": no Image or Rootfs key specified"

[Container]
