# tests/cases/userns.container
[Container]
Image=localhost/imagename
Exec=top
UserNS=keep-id:uid=200,gid=210
## assert-podman-args "--userns" "keep-id:uid=200,gid=210"